# systemd/media-mover.timer
[Unit]
Description=Media Mover Timer

[Timer]
OnBootSec=5min
OnUnitActiveSec=10min
Persistent=true

[Install]
WantedBy=timers.target

